# templates/upstart.conf
### THIS FILE IS MANAGED BY PUPPET
### DO NOT MODIFY IF YOU KNOW WHAT'S GOOD FOR YOU

description "tomcat <%= name %>"
author "puppet"

start on stopped networking
stop on shutdown

setuid <%= owner %>

env JAVA_HOME=/opt/jdk16
env CATALINA_BASE=<%= catalina_base %>

exec /opt/tomcat/bin/catalina.sh run
